# supervisord.conf
[unix_http_server]
file=/run/supervisor.sock

[supervisord]
nodaemon=true
# root required: audio device access (/dev/snd), SYS_BOOT for reboot endpoint
user=root
logfile=/dev/null
loglevel=info

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[program:webui]
command=python3 /app/app.py
directory=/app
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=5

[program:init]
; Reads config.json once on startup and starts the enabled audio services.
command=/usr/local/bin/init-services.sh
autostart=true
autorestart=false
exitcodes=0
startretries=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:librespot]
; Started/stopped by init-services.sh and the web UI — not by supervisord directly.
; Zeroconf mode: device appears in the Spotify app by hostname, no credentials needed.
command=/usr/local/bin/run-librespot.sh
autostart=false
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[program:snapclient]
command=/usr/local/bin/run-snapclient.sh
autostart=false
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30
